# deploy/exo-server.service
[Unit]
Description=Unified Server FastAPI backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nanotech
Group=nanotech
WorkingDirectory=/opt/unified_exo/server
EnvironmentFile=/etc/exo-platform/server.env
ExecStart=/opt/unified_exo/server/.venv/bin/python /opt/unified_exo/server/app.py
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/lib/exo-platform /opt/unified_exo/server

[Install]
WantedBy=multi-user.target
