# deploy/unified-server-frontend.service
[Unit]
Description=Unified Server — Next.js frontend
After=network-online.target unified-server.service
Wants=network-online.target

[Service]
Type=simple
User=nanotech
Group=nanotech
WorkingDirectory=/opt/unified_server/nanotech_website
EnvironmentFile=/etc/unified-server/frontend.env
Environment=NODE_ENV=production
Environment=HOSTNAME=0.0.0.0
Environment=PORT=8080
Environment=NEXT_TELEMETRY_DISABLED=1
ExecStart=/usr/bin/pnpm exec next start -p 8080 -H 0.0.0.0
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/unified_server/nanotech_website/.next /opt/unified_server/nanotech_website/node_modules

[Install]
WantedBy=multi-user.target
