# eeg_lab/deploy/eeg-lab.service
[Unit]
Description=EEG Lab Server
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/eeg-lab
ExecStart=/opt/eeg-lab/venv/bin/uvicorn backend.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
